# psql-ha-manual/configs/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── etcd ──────────────────────────────────────────────────────────────────────
[program:etcd]
command=/usr/local/bin/etcd
    --name etcd1
    --data-dir /var/lib/etcd
    --initial-advertise-peer-urls http://172.21.0.11:2380
    --listen-peer-urls            http://0.0.0.0:2380
    --listen-client-urls          http://0.0.0.0:2379
    --advertise-client-urls       http://172.21.0.11:2379
    --initial-cluster-token       pg-ha-cluster
    --initial-cluster-state       new
    --initial-cluster             etcd1=http://172.21.0.11:2380,etcd2=http://172.21.0.12:2380,etcd3=http://172.21.0.13:2380
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/etcd.log
stderr_logfile=/var/log/supervisor/etcd.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3

# ── patroni ───────────────────────────────────────────────────────────────────
[program:patroni]
command=/opt/patroni/bin/patroni /etc/patroni/patroni.yml
autostart=false
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/patroni.log
stderr_logfile=/var/log/supervisor/patroni.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
